# templates/iptv-weight-agent@.service.j2
[Unit]
Description=IPTV Weight Agent for port %i
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/iptv-weight-agent/%i
ExecStart=/usr/local/bin/iptv-weight-agent %i {{ role_haproxy__agent_base_port }} "$PROVIDER_URL" "$XTREAM_USER" "$XTREAM_PASS"
Restart=always
RestartSec=2
# Keep secrets out of logs
LogLevelMax=notice

[Install]
WantedBy=multi-user.target
